# backend/tfw/base/solvable/supervisor/event_handler_main.conf
[program:event_handler_main]
user=root
directory=%(ENV_TFW_EHMAIN_DIR)s
command=python3 -u event_handler_main.py
priority=3

[supervisord]
strip_ansi=false
