# backend/avatao_startr/tfw/base/solvable/supervisor/pipe_io_main.conf
[program:pipe_io_main]
user=root
directory=%(ENV_TFW_EHMAIN_DIR)s
command=python3 -u pipe_io_main.py
priority=4
